# script/sheepdog.service.in
[Unit]
Description=Sheepdog QEMU/KVM Block Storage
After=network.target
Wants=syslog.target

[Service]
EnvironmentFile=-@SYSCONFDIR@/conf.d/sheepdog.conf
EnvironmentFile=-@SYSCONFDIR@/sysconfig/sheepdog
EnvironmentFile=-@SYSCONFDIR@/default/sheepdog
ExecStart=@SBINDIR@/sheep --pidfile @LOCALSTATEDIR@/run/sheep.pid $SHEEP_OPTS
PIDFile=@LOCALSTATEDIR@/run/sheep.pid
Type=forking

[Install]
WantedBy=multi-user.target
