# script/sheepdog.service.in
[Unit]
Description=Sheepdog QEMU/KVM Block Storage
After=network.target
Wants=syslog.target

[Service]
EnvironmentFile=-/etc/conf.d/sheepdog.conf
EnvironmentFile=-/etc/sysconfig/sheepdog
EnvironmentFile=-/etc/default/sheepdog
ExecStart=/usr/sbin/sheep --pidfile /var/run/sheep.pid --log dst=syslog $SHEEP_OPTS /var/lib/sheepdog
PIDFile=/var/run/sheep.pid
Type=forking

[Install]
WantedBy=multi-user.target
